# systemd/telecodex.service
[Unit]
Description=TeleCodex Telegram bridge
Documentation=file:%h/telecodex/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/telecodex
Environment=NODE_ENV=production
EnvironmentFile=%h/telecodex/.env
ExecStart=/usr/bin/env TELECODEX_INSTANCE=default node dist/index.js
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=20

[Install]
WantedBy=default.target
