# systemd/telecodex@.service
[Unit]
Description=TeleCodex Telegram bridge (%i)
Documentation=file:%h/telecodex/README.md
After=network-online.target telecodex-codex-app-server.service
Wants=network-online.target telecodex-codex-app-server.service

[Service]
Type=simple
WorkingDirectory=%h/telecodex
Environment=NODE_ENV=production
EnvironmentFile=%h/telecodex/.env.%i
ExecStart=/usr/bin/env TELECODEX_INSTANCE=%i node dist/index.js
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=20

[Install]
WantedBy=default.target
