# systemd/telecodex.service
[Unit]
Description=TeleCodex Telegram bridge
Documentation=file:%h/telecodex/README.md
After=network-online.target telecodex-codex-app-server.service
Wants=network-online.target telecodex-codex-app-server.service

[Service]
Type=simple
WorkingDirectory=%h/telecodex
Environment=NODE_ENV=production
Environment="NODE_OPTIONS=--dns-result-order=ipv4first --no-network-family-autoselection"
EnvironmentFile=%h/telecodex/.env
ExecStart=/usr/bin/env TELECODEX_INSTANCE=default node dist/index.js
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=20

[Install]
WantedBy=default.target
